# templates/etc/profile.d/alias.sh

alias sudo='sudo '
alias apt-get='echo "This machine is under Ansible control - please use a playbook to modify anything"'
